# deploy/funasr-nano-asr.conf
[program:funasr-nano-asr]
directory=/root/autodl-tmp/funasr-nano-service
command=/bin/bash -lc 'source /root/autodl-tmp/funasr-nano-venv/bin/activate; export MODELSCOPE_CACHE=/root/autodl-tmp/models; export HF_HOME=/root/autodl-tmp/hf-home; export TRANSFORMERS_CACHE=/root/autodl-tmp/hf-home/transformers; export CUDA_VISIBLE_DEVICES=0; exec python -u serve_realtime_ws_phdebate.py --host 127.0.0.1 --port 10095 --model FunAudioLLM/Fun-ASR-Nano-2512 --hub ms --device cuda:0 --language 中文 --dtype bf16 --gpu-memory-utilization 0.35 --max-model-len 2048 --decode-interval 0.48 --disable-spk'
autostart=true
autorestart=true
startsecs=20
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/root/autodl-tmp/logs/funasr-nano-asr.log
stderr_logfile=/root/autodl-tmp/logs/funasr-nano-asr.err.log
environment=PYTHONUNBUFFERED="1"
