# deploy/phdebate-stack.supervisor.conf
[program:phdebate]
priority=10
directory=/root/autodl-tmp/phdebate/apps/backend
command=/bin/bash -lc 'set -a; source /root/autodl-tmp/phdebate/.env; set +a; exec /root/autodl-tmp/phdebate/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 6006 --no-access-log'
environment=PYTHONPATH="/root/autodl-tmp/phdebate/apps/backend",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=5
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/root/autodl-tmp/phdebate/logs/phdebate.stdout.log
stderr_logfile=/root/autodl-tmp/phdebate/logs/phdebate.stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:phdebate-voice-agent]
priority=20
directory=/root/autodl-tmp/phdebate/apps/voice_agent
command=/bin/bash -lc 'set -a; source /root/autodl-tmp/phdebate/.env; set +a; exec /root/autodl-tmp/phdebate/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 6008 --no-access-log'
environment=PYTHONPATH="/root/autodl-tmp/phdebate/apps/voice_agent",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=5
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/root/autodl-tmp/phdebate/logs/voice-agent.stdout.log
stderr_logfile=/root/autodl-tmp/phdebate/logs/voice-agent.stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:qwen3-asr]
priority=30
directory=/root/autodl-tmp/qwen3-asr-openai
command=/bin/bash -lc 'set -a; source /root/autodl-tmp/qwen3-asr-openai/.env; set +a; exec /root/autodl-tmp/qwen3-asr-openai/venv/bin/python -u server.py'
environment=PYTHONUNBUFFERED="1"
autostart=false
autorestart=true
startsecs=20
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/root/autodl-tmp/qwen3-asr-openai/logs/asr.stdout.log
stderr_logfile=/root/autodl-tmp/qwen3-asr-openai/logs/asr.stderr.log
stdout_logfile_maxbytes=100MB
stderr_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:qwen3-tts]
priority=40
directory=/root/autodl-tmp/qwen3-tts-openai
command=/root/autodl-tmp/qwen3-tts-openai/start_after_asr.sh
environment=PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=20
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/root/autodl-tmp/qwen3-tts-openai/logs/tts.stdout.log
stderr_logfile=/root/autodl-tmp/qwen3-tts-openai/logs/tts.stderr.log
stdout_logfile_maxbytes=100MB
stderr_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile_backups=5
