# deploy/phdebate-supervisord.conf
[program:phdebate]
directory=/root/autodl-tmp/phdebate/apps/backend
command=/root/autodl-tmp/phdebate/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 6006
environment=PYTHONPATH="/root/autodl-tmp/phdebate/apps/backend"
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/root/autodl-tmp/phdebate/logs/phdebate.stdout.log
stderr_logfile=/root/autodl-tmp/phdebate/logs/phdebate.stderr.log

[program:phdebate-voice-agent]
directory=/root/autodl-tmp/phdebate/apps/voice_agent
command=/root/autodl-tmp/phdebate/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 6008
environment=PYTHONPATH="/root/autodl-tmp/phdebate/apps/voice_agent"
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/root/autodl-tmp/phdebate/logs/voice-agent.stdout.log
stderr_logfile=/root/autodl-tmp/phdebate/logs/voice-agent.stderr.log

